# raspberryfluke.service
[Unit]
Description=RaspberryFluke LLDP/CDP E-Paper Display
After=network.target lldpd.service
Wants=lldpd.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/raspberryfluke
Environment=PYTHONPATH=/opt/raspberryfluke PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/raspberryfluke/raspberryfluke.py
Restart=on-failure
RestartSec=3

# Systemd hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true
RestrictSUIDSGID=true
RemoveIPC=true
RestrictAddressFamilies=AF_UNIX AF_NETLINK

[Install]
WantedBy=multi-user.target
